# deployment/shellscript.service
[Unit]
Description=My Custom Shell Script

[Service]
ExecStart=/tmp/shellscript.sh

[Install]
WantedBy=multi-user.target