# service_centos7/shellscript.service
[Unit]
Description=My Custom Shell Script

[Service]
ExecStart=/home/shellscript.sh

[Install]
WantedBy=multi-user.target